# init/permaweb-host.service
[Unit]
Description=Permaweb Host

[Service]
ExecStart=/root/go/bin/permaweb-host
Restart=on-failure

[Install]
WantedBy=default.target
